# intervisions.service
[Unit]
Description=InterVisions Annotation Platform
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/intervisions_annotation
Environment="DATABASE_PATH=/data/intervisions/intervisions.db"
Environment="UPLOAD_FOLDER=/data/intervisions/images"
Environment="SECRET_KEY=9e481407de6a6063f84cf4cf5a8ae6a627762c792206ecc02f472b2732a7baeb"
Environment="ADMIN_PASSWORD=1nt3rvisions2026"
ExecStart=/home/ubuntu/intervisions_annotation/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 2 --timeout 120 app.main:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
